# alloy_playbooks/software-exporter.timer
[Unit]
Description=Run software exporter hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
